# deploy/db/systemd/adsecute-db-healthcheck.service
[Unit]
Description=Adsecute PostgreSQL Health Check
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/adsecute-db-healthcheck.sh
